# systemd/.config/systemd/user/battery-warning.service
# battery-warning.service

[Unit]
Description=Battery warning

[Service]
Type=simple
ExecStart=/usr/bin/env %h/bin/i3-battery-warning
