# systemd/.config/systemd/user/mbsync.service
# mbsync.service

[Unit]
Description=Mbsync Quicksync

[Service]
Type=oneshot
ExecStart=/usr/bin/env %h/bin/mail-sync
